# autowakeup.timer
[Unit]
Description=Autowakeup timer

[Timer]
OnBootSec=15s
OnUnitActiveSec=1h 

[Install]
WantedBy=timers.target
